# deploy/monitor-catchup.service
[Unit]
Description=general_monitor gap catch-up (short idempotent backfill to heal brief ingest outages)
After=network.target

[Service]
Type=oneshot
User=monitor
Group=monitor
WorkingDirectory=/opt/monitor
EnvironmentFile=/opt/monitor/.env
ExecStart=/opt/monitor/.venv/bin/python -m scripts.backfill all --days 2 --hours 48

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ReadWritePaths=/var/lib/monitor
LockPersonality=true
# MemoryDenyWriteExecute is intentionally NOT set: polars calls mprotect(PROT_WRITE|PROT_EXEC)
# at import time for its CPU feature check, which this flag would block.
