# deploy/monitor-ingest.timer
[Unit]
Description=Run general_monitor ingester every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
